# 49-htcpnet.rules
SUBSYSTEM=="usb", ATTR{idVendor}=="0bb2", ATTR{idProduct}=="0fb4",ACTION=="add", RUN+="/usr/bin/sudo -b /usr/local/bin/htcnet.sh"
